# scripts/io.restorationai.autoblog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.restorationai.autoblog</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>/Users/santino/Desktop/mywebsitecode/Restoration-AI-website-2026-main/scripts/auto-blog.sh</string>
  </array>

  <!-- Mon (Weekday 1) and Thu (Weekday 4) at 09:00 local time -->
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key><integer>1</integer>
      <key>Hour</key><integer>9</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>4</integer>
      <key>Hour</key><integer>9</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
  </array>

  <!-- If the Mac was asleep at 9am, run as soon as it wakes. -->
  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/santino/Desktop/mywebsitecode/Restoration-AI-website-2026-main/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/santino/Desktop/mywebsitecode/Restoration-AI-website-2026-main/logs/launchd.err.log</string>
</dict>
</plist>
